# ansible/etc/profile.d/kylin.env.sh
export KYLIN_HOME=/opt/kylin
export PATH=$KYLIN_HOME/bin:$PATH
